# root_files/defaults/etc/hotplug.d/iface/10-ethtool
[ ifup = "$ACTION" ] || exit 0
[ -n "$DEVICE" ] || exit 0

if [ "$DEVICE" == wan ] ; then
      ethtool --set-eee "$DEVICE" eee off
fi
