# root_files/etc/hotplug.d/iface/50-autoipd
# create avahi state directory
mkdir -p /var/lib

[ ifup = "$ACTION" ] && {
	[ -n "$DEVICE" ] && {
		AUTOIP_OPT=$(uci -q get network.${INTERFACE}.autoip)
		[ "${AUTOIP_OPT}" == "yes" ] && {
			avahi-autoipd -D ${DEVICE}
		}
	}
}

[ ifdown = "$ACTION" ] && {
	[ -n "$DEVICE" ] && {
		AUTOIP_OPT=$(uci -q get network.${INTERFACE}.autoip)
		[ "${AUTOIP_OPT}" == "yes" ] && {
			avahi-autoipd -k ${DEVICE}
		}
	}
}
