# root_files/basilicata/etc/uci-defaults/99-zerotier
#!/bin/sh

uci -q batch <<-EOF >/dev/null
  set network.ZeroTier=interface
  set network.ZeroTier.proto='none'
  set network.ZeroTier.device='owzt9dd842'
  commit network
EOF
