# root_files/etc/uci-defaults/99-pppoe-wan-dnsmasq
#!/bin/sh

uci add_list dhcp.@dnsmasq[0].notinterface='pppoe-wan'
uci commit dhcp
